# usr/local/bin/setzone.service
[Unit]
	Description=setzone.service will set the timezone and locale after boot

	[Service]
	Type=simple
	ExecStart=/usr/local/bin/settimezone

	[Install]
	WantedBy=multi-user.target
